# com.workflow-distiller.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 模板：占位符由 install.sh 在安装时填充（__PROJECT_ROOT__ / __PYTHON__ / __HOME__）。
     勿直接 load 本模板；load 的是 install.sh 生成到 ~/Library/LaunchAgents 的实例。 -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.workflow-distiller</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>distiller.retro</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PYTHONPATH</key>
        <string>__PROJECT_ROOT__</string>
    </dict>

    <!-- 每周一 09:00 触发每周复盘 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__PROJECT_ROOT__/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_ROOT__/logs/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
